# meta-ai-assistant/recipes-camera/camera-stream/files/camera-stream.service
[Unit]
Description=Camera MJPEG Stream Server
After=network.target

[Service]
Type=simple
ExecStart=/opt/camera/camera-stream
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
